# front_end/system_files/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true

[program:gunicorn]
command=/venv/bin/gunicorn --bind 0.0.0.0:5000 --certfile=/etc/nginx/ssl/self-signed.crt --keyfile=/etc/nginx/ssl/self-signed.key --timeout 120 --workers 9 --worker-class gevent app:app
autorestart=true
directory=/app
stderr_logfile=/app/logs/gunicorn.err.log
stdout_logfile=/app/logs/gunicorn.out.log