# rtmp/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/src/nginx/nginx-1.21.1/objs/nginx -g "daemon off;"

[program:flask]
command=python3 /app/app.py
