# rtmp/system/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/local/nginx/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true

[program:flask]
command=python3 /app/driver.py
